# provision/files/mt7.ini
[program:mt7]
user=vagrant
command=/etc/supervisord.d/mt7.sh
autostart=true
autorestart=true
stopsignal=QUIT
